# var/sense.service
[Unit]
Description=Weather Sensing

[Service]
ExecStart=/usr/local/bin/sense

[Install]
WantedBy=multi-user.target
